# root/etc/s6-overlay/s6-rc.d/init-hummingbot-chown/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ -z ${LSIO_READ_ONLY_FS} ]] && [[ -z ${LSIO_NON_ROOT_USER} ]]; then
    lsiown -R abc:abc /home/hummingbot
fi
